# usr/src/cmd/distro_const/slim_cd/slimcd_live.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    CDDL HEADER START
   
    The contents of this file are subject to the terms of the
    Common Development and Distribution License (the "License").
    You may not use this file except in compliance with the License.
   
    You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
    or http://www.opensolaris.org/os/licensing.
    See the License for the specific language governing permissions
    and limitations under the License.
   
    When distributing Covered Code, include this CDDL HEADER in each
    file and include the License file at usr/src/OPENSOLARIS.LICENSE.
    If applicable, add the following below this CDDL HEADER, with the
    fields enclosed by brackets "[]" replaced with your own identifying
    information: Portions Copyright [yyyy] [name of copyright owner]
   
    CDDL HEADER END
   
    Copyright 2010 Sun Microsystems, Inc.  All rights reserved.
    Use is subject to license terms.

    Service profile customization for Live CD Image.
-->
<service_bundle type='profile' name='live_cd'
	 xmlns:xi='http://www.w3.org/2003/XInclude' >

  <!--
      svc.startd(1M) services
  -->

  <service name='system/hal' version='1' type='service'>
    <instance name='default' enabled='true'/>
  </service>
  <service name='system/postrun' version='1' type='service'>
    <instance name='default' enabled='false'/>
  </service>
  <service name='system/filesystem/zfssnap-roleadd' version='1' type='service'>
    <instance name='default' enabled='false'/>
  </service>
  <service name='network/nfs/server' version='1' type='service'>
    <instance name='default' enabled='false'/>
  </service>
  <service name='network/physical' version='1' type='service'>
      <instance name='default' enabled='false'/>
      <instance name='nwam' enabled='true'/>
  </service>
  <service name='network/rpc/gss' version='1' type='service'>
      <instance name='default' enabled='false'/>
  </service>
  <service name='network/sendmail-client' version='1' type='service'>
    <instance name='default' enabled='false'/>
  </service>
  <service name='network/shares/group' version='1' type='service'>
    <instance name='default' enabled='false'/>
    <instance name='zfs' enabled='false'/>
  </service>
  <service name='network/tnctl' version='1' type='service'>
    <instance name='default' enabled='false'/>
  </service>
  <service name='network/inetd-ugrade' version='1' type='service'>
    <instance name='default' enabled='false'/>
  </service>
  <service name='system/filesystem/rmvolmgr' version='1' type='service'>
    <instance name='default' enabled='true'/>
  </service>
  <service name='application/management/sma' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='application/management/seaport' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='application/management/snmpdx' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='application/management/wbem' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='application/print/cleanup' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='application/print/ipp-listener' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='application/print/ppd-cache-update' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='system/basicreg' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='platform/i86pc/kdmconfig' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='application/stosreg' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='network/security/ktkt_warn' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='network/dns/multicast' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='application/pkg/update' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='system/avahi-bridge-dsd' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='application/x11/xfs' version='1' type='service'>
    <instance name='default' enabled='true'/>
  </service>

  <!--
      Make sure GNOME cache services for the live image are all disabled.
  -->
  <service name='application/desktop-cache/desktop-mime-cache' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>

  <service name='application/desktop-cache/gconf-cache' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>

  <service name='application/desktop-cache/icon-cache' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>

  <service name='application/desktop-cache/input-method-cache' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>

  <service name='application/desktop-cache/mime-types-cache' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>

  <service name='application/desktop-cache/pixbuf-loaders-installer' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>

  <!--
      Increase the timeout to avoid error messages for the next 3 services.
      Default of 30 seconds is not enough when starting from a slow device
      like a DVD drive.
  -->
  <service name='system/dbus' version='1' type='service'>
    <instance name='default' enabled='true'>
       <exec_method
                type='method'
                name='start'
		exec='/lib/svc/method/svc-dbus start'
                timeout_seconds='90'/>
    </instance>
  </service>

  <!--
      Enable gdm services.
  -->
  <service name='application/graphical-login/gdm' version='1' type='service'>
    <instance name='default' enabled='true'>
       <exec_method
                type='method'
                name='start'
		exec='/lib/svc/method/svc-gdm start'
                timeout_seconds='120'/>
    </instance>
  </service>

  <!--
      Enable X11 services
  -->
  <service name='application/opengl/ogl-select' version='1' type='service'>
    <instance name='default' enabled='true'>
       <exec_method
                type='method'
                name='start'
		exec='/lib/svc/method/ogl-select start'
                timeout_seconds='90'/>
    </instance>
  </service>
</service_bundle>
