# usr/src/cmd/system-config/svc/system-config.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

 Copyright (c) 2010, 2011, Oracle and/or its affiliates. All rights reserved.

 NOTE:  This service manifest is not editable; its contents will
 be overwritten by package or patch operations, including
 operating system upgrade.  Make customizations in a different
 file.
-->
<service_bundle type='manifest' name='system:system-config'>

<service
	name='system/config'
	type='service'
	version='1'>

	<create_default_instance enabled='false'/>

	<single_instance/>

	<!-- system/config should depend upon system/filesystem/root
	     since it works with following assumptions which in general
	     are not guaranteed until system/filesystem/root comes online:
	      - writable root filesystem
	      - mounted /usr filesystem

	     However, at this point system/config can't depend upon
	     system/filesystem/root, since that would lead to following cyclic
	     dependency:
	     system/config -> system/filesystem/root -> system/metainit ->
	     system/identity:node -> system/config

	     The problem is tracked by CR7011249.

	     That is not a big deal in standard running environment with
	     root filesystem on ZFS (which is mounted as writable by
	     kernel) and /usr not created on a separate ZFS dataset (thus it is
	     available once root filesystem is mounted). -->
<!--
	<dependency name='fs'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri
		    value='svc:/system/filesystem/root:default' />
	</dependency>
-->

	<!-- system/device/local takes care of populating device tree
	     during first post-install boot of global zone. This is one
	     of prerequisites which dladm(1m) depends on (in global zone only)
	     in order to work correctly. dladm(1m) is consumed by System
	     Configuration Interactive tool for obtaining information about
	     network interfaces.

	     However, at this point system/config can't depend upon
	     system/device/local, since that would lead to following cyclic
	     dependency:
	     system/config -> system/device/local -> system/filesystem/usr ->
	     system/boot-archive -> system/filesystem/root -> system/metainit ->
	     system/identity:node -> system/config

	     Again, this is the problem tracked by CR7011249.

	     As a stop-gap solution, system/config start method calls
	     devfsadm(1m) in case SCI tool is to be invoked in global zone. -->
<!--
	<dependency name='device-local'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri
		    value='svc:/system/device/local:default' />
	</dependency>
-->

	<!-- Dependency on network/datalink-management guarantees that network
             interfaces eligible for configuration can be discovered by System
             Configuration Interactive tool. -->
	<dependency name='datalink-management'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri
		    value='svc:/network/datalink-management' />
	</dependency>

	<dependent
	    name='systemconfig_singleuser'
            grouping='optional_all'
            restart_on='none'>
            <service_fmri value='svc:/milestone/single-user' />
	</dependent>

	<dependent
	    name='systemconfig_manifest_import'
            grouping='optional_all'
            restart_on='none'>
            <service_fmri value='svc:/system/manifest-import:default' />
	</dependent>

	<!-- Following services can be configured by System Configuration
	     Interactive (SCI) tool invoked from this (system/config) smf
	     service. Thus make sure those services are not started before
	     SCI tool finishes its job. -->
	<dependent
	    name='systemconfig_network_install'
            grouping='optional_all'
            restart_on='none'>
            <service_fmri value='svc:/network/install' />
	</dependent>

	<dependent
	    name='systemconfig_network_physical_default'
            grouping='optional_all'
            restart_on='none'>
            <service_fmri value='svc:/network/physical:default' />
	</dependent>

	<dependent
	    name='systemconfig_network_dns_client'
            grouping='optional_all'
            restart_on='none'>
            <service_fmri value='svc:/network/dns/client' />
	</dependent>

	<dependent
	    name='systemconfig_system_identity_node'
            grouping='optional_all'
            restart_on='none'>
            <service_fmri value='svc:/system/identity:node' />
	</dependent>

	<dependent
	    name='systemconfig_system_keymap'
            grouping='optional_all'
            restart_on='none'>
            <service_fmri value='svc:/system/keymap' />
	</dependent>

	<dependent
	    name='systemconfig_console_login'
            grouping='optional_all'
            restart_on='none'>
            <service_fmri value='svc:/system/console-login' />
	</dependent>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/svc-system-config'
		timeout_seconds='0' />

	<exec_method
		type='method'
		name='stop'
		exec=':true'
		timeout_seconds='0' />

	<exec_method
		type='method'
		name='refresh'
		exec=':true'
		timeout_seconds='0' />

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring'
			value='transient' />
	</property_group>

	<!-- The following property group is used only at install
	     time to configure an initial user account -->
	<property_group name='user_account' type='application'>
		<stability value='Evolving'/>
		<propval name='read_authorization' type='astring'
		    value='solaris.smf.read.system-config' />
		<propval name='login' type='astring' value=''/>
		<propval name='password' type='astring' value=''/>
		<propval name='description' type='astring' value=''/>
		<propval name='shell' type='astring' value=''/>
		<propval name='uid' type='count' value='0'/>
		<propval name='gid' type='count' value='0'/>
		<propval name='profiles' type='astring' value=''/>
		<propval name='roles' type='astring' value=''/>
		<propval name='sudoers' type='astring' value=''/>
		<propval name='type' type='astring' value=''/>
		<propval name='expire' type='astring' value=''/>
		<propval name='home_zfs_dataset' type='astring' value=''/>
		<propval name='home_mountpoint' type='astring' value=''/>
		<propval name='autohome' type='astring' value=''/>
	</property_group>

	<!-- The following property group is used only at install
	     time to configure an initial root account --> 
	<property_group name='root_account' type='application'>
		<stability value='Evolving' />
		<propval name='read_authorization' type='astring'
		    value='solaris.smf.read.system-config' />
		<propval name='password' type='astring' value='' />
		<propval name='type' type='astring' value=''/>
		<propval name='expire' type='astring' value=''/>
	</property_group>

	<!-- The following property group is used to control
	     process of system configuration/unconfiguration --> 
	<property_group name='configuration' type='application'>
		<stability value='Evolving' />

		<!-- Setting following property to 'true' triggers overall
		     process of system configuration. If set to 'false',
		     'interactive_config' and 'groupings' properties have
		     no effect.
		     It is reset to 'false' after and if system configuration
		     succeeded. -->
		<propval name='configure' type='boolean' value='false' />

		<!-- If following property is set to 'true', SCI tool is
		     invoked during boot. It is reset to 'false' by
		     system/config smf service after SCI tool successfully
		     generates System Configuration profile. -->
		<propval name='interactive_config' type='boolean'
                    value='false' />

		<!-- Following property specifies which portions of system
		     should be configured. It could contain combination
		     of following values (separated by spaces):

		      - system - everything is unconfigured (used as a default)
		      - network
		      - naming_services
		      - location (locale, keyboard layout, timezone)
		      - users (user and root account)
		      - identity (system hostname)
		-->
		<propval name='config_groups' type='astring' value='system'/>
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			system configuration
			</loctext>
		</common_name>
	</template>
</service>
</service_bundle>

