# usr/src/cmd/distro_const/profile/ai.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    CDDL HEADER START
   
    The contents of this file are subject to the terms of the
    Common Development and Distribution License (the "License").
    You may not use this file except in compliance with the License.
   
    You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
    or http://www.opensolaris.org/os/licensing.
    See the License for the specific language governing permissions
    and limitations under the License.
   
    When distributing Covered Code, include this CDDL HEADER in each
    file and include the License file at usr/src/OPENSOLARIS.LICENSE.
    If applicable, add the following below this CDDL HEADER, with the
    fields enclosed by brackets "[]" replaced with your own identifying
    information: Portions Copyright [yyyy] [name of copyright owner]
   
    CDDL HEADER END
   
    Copyright (c) 2008, 2011, Oracle and/or its affiliates. All rights reserved.

    Service profile customization for Automated Installer image
-->
<service_bundle type='profile' name='ai'>
  <!--
      svc.startd(1M) services
  -->
  <service name='network/physical' version='1' type='service'>
      <instance name='default' enabled='true'/>
      <instance name='nwam' enabled='false'/>
  </service>
  <service name='network/routing/ndp' version='1' type='service'>
    <instance name='default' enabled='false'/>
  </service>
  <service name='network/shares/group' version='1' type='service'>
    <instance name='default' enabled='false'/>
  </service>
  <service name='network/dns/multicast' version='1' type='service'>
    <instance name='default' enabled='true'/>
  </service>
  <service name='application/manifest-locator' version='1' type='service'>
    <instance name='default' enabled='true' />
  </service>
  <service name='application/auto-installer' version='1' type='service'>
    <instance name='default' enabled='true' />
  </service>

  <service name='system/name-service/switch' version='1' type='service'>
    <exec_method
      type='method'
      name='start'
      timeout_seconds='0'
      exec=':true'/>
    <exec_method
      type='method'
      name='refresh'
      timeout_seconds='0'
      exec=':true'/>
  </service> 

</service_bundle>
