# usr/src/cmd/auto-install/svc/auto-installer.xml
<?xml version="1.0"?>
<!--
 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

 Copyright (c) 2009, 2011, Oracle and/or its affiliates. All rights reserved.

 NOTE:  This service manifest is not editable; its contents will
 be overwritten by package or patch operations, including
 operating system upgrade.  Make customizations in a different
 file.
-->

<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='SUNWauto-install:auto-installer'>

<service
	name='application/auto-installer'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />
	<single_instance/>

	<!-- There is a soft dependency on console-login. -->
	<dependency
		name='console-login'
		grouping='optional_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/system/console-login' />
	</dependency>

        <!-- Do not start the auto-installer until devchassis daemon 
             SMF service is online to ensure target selection works 
             correctly for manifests that specifies receptacle
             information for the drives -->
        <dependency
                name='devchassis-daemon'
                grouping='require_all'
                restart_on='none'
                type='service'>
                <service_fmri value='svc:/system/devchassis:daemon' />
        </dependency>

	<dependency
		name='manifest-locator'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/application/manifest-locator' />
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/auto-installer'
		timeout_seconds='0' />

	<exec_method
		type='method'
		name='stop'
		exec=':true'
		timeout_seconds='0' />

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient' />
		<propval name='ignore_error' type='astring' value='core,signal' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				automated installer
			</loctext>
		</common_name>
		<documentation>
			<manpage title='auto-install' section='1M'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>
</service_bundle>
